# sdcard/sdcard.rules
SUBSYSTEM=="block",ACTION=="add",RUN+="/home/sdcard-autorun"
#RUN+="/home/sdcard-autorun"
